# ops/systemd/ooc-find-gigs.timer
[Unit]
Description=Trigger OOC find-gigs weekly, Sun 05:00 local
Requires=ooc-find-gigs.service

[Timer]
OnCalendar=Sun *-*-* 05:00:00
RandomizedDelaySec=30min
Persistent=true
Unit=ooc-find-gigs.service

[Install]
WantedBy=timers.target
