# ops/systemd/ooc-verify-gigs.timer
[Unit]
Description=Trigger OOC verify-gigs weekly, Sun 02:00 local
Requires=ooc-verify-gigs.service

[Timer]
OnCalendar=Sun *-*-* 02:00:00
RandomizedDelaySec=30min
Persistent=true
Unit=ooc-verify-gigs.service

[Install]
WantedBy=timers.target
